# conf.d/xvfb.conf
[program:xvfb]
command=Xvfb :0 -screen 0 "1400"x"900"x24
autorestart=true
