# conf.d/xvfb.conf
[program:xvfb]
command=Xvfb :1 -screen 1 "1400"x"600"x24
autorestart=true
